# com.binance.scalper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.binance.scalper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/rahulpatil/Desktop/binance_scalper/start_bot.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rahulpatil/Desktop/binance_scalper</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/rahulpatil/Desktop/binance_scalper/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rahulpatil/Desktop/binance_scalper/logs/launchd_stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
